# module/local/crontab
@reboot sleep 300 && ${game_dir}/backup_sunkenland.sh &>> ${game_dir}/backup.log
@hourly ${game_dir}/backup_sunkenland.sh &>> ${game_dir}/backup.log
